# deploy/sagatoy.service
[Unit]
Description=SagaToy/ToToyAI FastAPI Backend
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=notify
User=harvard
Group=harvard
WorkingDirectory=/var/www/sagatoy

# Environment variables
Environment="PATH=/var/www/sagatoy/venv/bin"
EnvironmentFile=/var/www/sagatoy/.env

# Start command
ExecStart=/var/www/sagatoy/venv/bin/uvicorn totoyai.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info \
    --access-log \
    --use-colors

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/sagatoy

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sagatoy

[Install]
WantedBy=multi-user.target
